# .config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Service]
Environment="XDG_CURRENT_DESKTOP=wlr"
Environment="WAYLAND_DISPLAY=wayland-1"
