# .config/systemd/user/waydroid-monitor.service.d/override.conf
[Service]
Environment="SOURCE=/data/main/projectpak"
Environment="TARGET=/home/chent/.local/share/waydroid/data/media/0/Documents/"